# automation/com.ai-edu-system.autosync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ai-edu-system.autosync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/peterpan528/ai-edu-system/automation/auto-sync.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/peterpan528/Library/Logs/ai-edu-system-autosync.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/peterpan528/Library/Logs/ai-edu-system-autosync.log</string>
</dict>
</plist>
